# service/species-qa.service
[Unit]
Description=Metadata QA service (species + locations)
After=network.target

[Service]
Type=simple
User=egaillac
WorkingDirectory=/home/egaillac/MetaP/species_qa_service
Environment=SPECIES_QA_OTT_DB=/home/egaillac/MetaP/classifier/data/processed/ott_index.sqlite
Environment=SPECIES_QA_MODEL_DIR=/home/egaillac/MetaP/qa_bert/models/boe-qa-v2
Environment=TOKENIZERS_PARALLELISM=false
Environment=HOST=0.0.0.0
Environment=PORT=8010
ExecStart=/home/egaillac/MetaP/species_qa_service/start.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
